# config/docker/php-fpm/workers.ini
[program:worker.product_finder]
command=php /var/www/html/bin/console swarrot:consume:product_finder import.finder
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/var/www/html
autostart=true
autorestart=true
startsecs=1
startretries=100
redirect_stderr=false
stdout_capture_maxbytes=10MB
stderr_capture_maxbytes=10MB

[program:worker.product_mapping_field]
command=php /var/www/html/bin/console swarrot:consume:product_mapping_field product.mapping_field
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/var/www/html
autostart=true
autorestart=true
startsecs=1
startretries=100
redirect_stderr=false
stdout_capture_maxbytes=10MB
stderr_capture_maxbytes=10MB

[program:worker.product_mapping_value]
command=php /var/www/html/bin/console swarrot:consume:product_mapping_value product.mapping_value
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/var/www/html
autostart=true
autorestart=true
startsecs=1
startretries=100
redirect_stderr=false
stdout_capture_maxbytes=10MB
stderr_capture_maxbytes=10MB

[program:worker.product_import]
command=php /var/www/html/bin/console swarrot:consume:product_import import.product
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/var/www/html
autostart=true
autorestart=true
startsecs=1
startretries=100
redirect_stderr=false
stdout_capture_maxbytes=10MB
stderr_capture_maxbytes=10MB

[program:worker.event_lang_upsert_translation]
command=php /var/www/html/bin/console swarrot:consume:event_lang_upsert_translation event.lang.upsert_translation
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/var/www/html
autostart=true
autorestart=true
startsecs=1
startretries=100
redirect_stderr=false
stdout_capture_maxbytes=10MB
stderr_capture_maxbytes=10MB

[program:worker.product_image_import]
command=php /var/www/html/bin/console swarrot:consume:product_image_import import.image
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/var/www/html
autostart=true
autorestart=true
startsecs=1
startretries=100
redirect_stderr=false
stdout_capture_maxbytes=10MB
stderr_capture_maxbytes=10MB